# ntopng/rootfs/etc/s6-overlay/s6-rc.d/netflow2ng/run
#!/command/with-contenv bashio
# ==============================================================================
# Li Tin O've Weedle Add-on: ntopng
# Runs the netflow2ng daemon
# ==============================================================================
declare -a options
declare netflow_port

bashio::log.info "Starting netflow2ng..."

options+=(--metrics="127.0.0.1:$(bashio::config 'netflow2ng_loc_port')")
options+=(--listen-zmq="tcp://127.0.0.1:$(bashio::config 'zmq_loc_port')")
options+=(--log-level="error")

netflow_port=$(bashio::addon.port '2055/udp')
if bashio::var.is_empty "${netflow_port}"; then
    bashio::log.warning
    bashio::log.warning "Netflow collector port is not defined"
    bashio::log.warning "Netflow will start only on localhost!"
    bashio::log.warning "This is probably not something you want"
    bashio::log.warning "as you won't be able to collect data from network!"
    bashio::log.warning
    options+=(--listen="127.0.0.1:2055")
else
    options+=(--listen="0.0.0.0:${netflow_port}")
fi

# Find the matching netflow log level
if bashio::config.has_value 'log_level'; then
    case "$(bashio::string.lower "$(bashio::config 'log_level')")" in
        all|trace)
            options+=(--log-level="trace")
            ;;
        debug)
            options+=(--log-level="debug")
            ;;
        info|notice)
            options+=(--log-level="info")
            ;;
        warning)
            options+=(--log-level="warn")
            ;;
        error|fatal|off)
            options+=(--log-level="error")
            ;;
    esac
fi

bashio::log.info "netflow2ng ${options[@]}"

exec s6-setuidgid netflow netflow2ng "${options[@]}"

