# ntopng/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv bashio
# ==============================================================================
# Li Tin O've Weedle Add-on: ntopng
# Runs the redis daemon
# ==============================================================================
declare -a options

bashio::log.info "Starting redis..."

options+=(--bind 127.0.0.1)
options+=(--port "$(bashio::config 'redis_loc_port')")
options+=(--daemonize no)
options+=(--dir /data/redis)

# Find the matching ntopng log level
if bashio::config.has_value 'log_level'; then
    case "$(bashio::string.lower "$(bashio::config 'log_level')")" in
        all|trace|debug)
            options+=(--loglevel debug)
            ;;
        info)
            options+=(--loglevel verbose)
            ;;
        notice)
            options+=(--loglevel notice)
            ;;
        warnng|error|fatal)
            options+=(--loglevel warning)
            ;;
        off)
            options+=(--loglevel warning --logfile "")
            ;;
    esac
fi

bashio::log.info "redis-server /etc/redis/redis.conf ${options[@]}"

exec s6-setuidgid redis redis-server /etc/redis/redis.conf "${options[@]}"

