# ntopng/rootfs/etc/s6-overlay/s6-rc.d/geoip/run
#!/command/with-contenv bashio
# ==============================================================================
# Li Tin O've Weedle Add-on: ntopng
# Runs the netflow2ng daemon
# ==============================================================================

while 1
do
    sleep 21600
    bashio::log.info "Updating Maxmind Database..."
    /usr/bin/geoipupdate
done
