# ntopng/rootfs/etc/s6-overlay/s6-rc.d/geoip/run
#!/command/with-contenv bashio
# ==============================================================================
# Li Tin O've Weedle Add-on: ntopng
# Runs the netflow2ng daemon
# ==============================================================================

bashio::log.info "Starting Maxmind Database update runner..."

while `true`
do
    sleep 21600
    if bashio::fs.file_exists /etc/GeoIP.conf; then
        rm -f /var/lib/GeoIP/.finished
        bashio::log.info "Updating Maxmind Database..."
        /usr/bin/geoipupdate
        bashio::log.info "Maxmind Database updated."
        touch /var/lib/GeoIP/.finished
    fi
done
