# ntopng/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Li Tin O've Weedle Add-on: ntopng
# Runs the Nginx daemon
# ==============================================================================

# Wait for ntopng to become available
bashio::net.wait_for "$(bashio::config 'ntopng_loc_port')" localhost 900

bashio::log.info "Starting NGinx..."
exec nginx
