# ntopng/rootfs/etc/s6-overlay/s6-rc.d/init-ntopng/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC2207
# ==============================================================================
# Li Tin O've Weedle Add-on: ntopng
# Handles configuration
# ==============================================================================

if ! bashio::fs.directory_exists /var/run/ntopng; then
    mkdir -p /var/run/ntopng
    chown -R ntopng:ntopng /var/run/ntopng
fi

if ! bashio::fs.directory_exists /data/rrd; then
    mkdir -p /data/rrd
    chown -R ntopng:ntopng /data/rrd
fi

if ! bashio::fs.directory_exists /data/pcap; then
    mkdir -p /data/pcap
    chown -R ntopng:ntopng /data/pcap
fi

if bashio::config.true 'custom_scripts'; then
    if ! bashio::fs.directory_exists /data/scripts; then
        mkdir -p /data/scripts
        cp -R /usr/local/share/ntopng/scripts /data/scripts
        chown -R ntopng:ntopng /data/scripts
    fi
fi

