# ntopng/rootfs/etc/s6-overlay/s6-rc.d/init-redis/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC2207
# ==============================================================================
# Home Assistant Community Add-on: ntopng
# Handles configuration
# ==============================================================================
readonly DATA="/data/ntopng"

if ! bashio::fs.directory_exists ${DATA}; then
    mkdir -p ${DATA}
    chmod 777 ${DATA}
    mkdir -p ${DATA}/redis
    chown -R ntopng:ntopng ${DATA}/redis
elif ! bashio::fs.directory_exists ${DATA}/redis; then
    mkdir -p ${DATA}/redis
    chown -R redis:redis ${DATA}/redis
fi
